# src/etc/cron.d/export-metrics
* * * * * root cd /opt/datadog_metrics/ && make stats >> /var/log/cron.log 2>&1
